# stockmonkey/openclaw/skills/stock_daily_brief/ai.stockmonkey.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>ai.stockmonkey.telegram-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/sameerhassen/Desktop/TECH CAREER/Mil by 30/StockMonkey/.venv/bin/python</string>
    <string>-m</string>
    <string>app.telegram_bot</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/sameerhassen/Desktop/TECH CAREER/Mil by 30/StockMonkey/stockmonkey</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/sameerhassen/Desktop/TECH CAREER/Mil by 30/StockMonkey/stockmonkey/logs/telegram_bot_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/sameerhassen/Desktop/TECH CAREER/Mil by 30/StockMonkey/stockmonkey/logs/telegram_bot_stderr.log</string>

</dict>
</plist>
